# Cassian-RV/aqualog-rv/web-portal/deploy/yxt-demo.service
[Unit]
Description=YangXiaTong temporary demo portal
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/dadatun-demo/current
ExecStart=/usr/bin/php -S 127.0.0.1:8081 -t /srv/dadatun-demo/current/public
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
